# deploy/quant-backtester.service
[Unit]
Description=Quant Backtester — A股量化回测 FastAPI 服务
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/.hermes/projects/quant-backtester
ExecStart=/usr/local/lib/hermes-agent/venv/bin/uvicorn quant_backtester.main:app --host 0.0.0.0 --port 8100 --log-level info
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1
MemoryMax=800M
CPUQuota=80%

[Install]
WantedBy=multi-user.target
